# 10777_monitoring.conf
description "monitoring"
author "stavros liaskos"

start on (started docker)
stop on runlevel [!2345]

# Automatically Respawn with finite limits
respawn
respawn limit 99 15
respawn limit unlimited

script
        # auto-build fails. But it's ok. We just want auto-container restart when fails
        # docker build -t xxxx/node_app /xxxxxx
        docker run -p 10777:3000 --name 10777_monitoring -d xxxxx/node_app
end script